# root/etc/uci-defaults/30_luci-theme-att
#!/bin/sh
uci batch <<-EOF
	set luci.themes.att=/luci-static/att
	set luci.main.mediaurlbase=/luci-static/att
	commit luci
EOF
exit 0
